# provisioning/roles/jumphost/files/etc/NetworkManager/dispatcher.d/30-promisc
#!/bin/bash

if [ "$1" == "eth1" ] && [ "$2" == "up" ]; then
    ip link set dev "$1" promisc on
fi
